# bin/afraid.timer
[Unit]
Description = Refresh freedns entry every hour
After = network.target

[Timer]
OnCalendar = hourly

[Install]
WantedBy=timers.target
