# rtc.service
[Unit]
Description=RTC script

[Service]
ExecStart=/store/scripts/rtc

[Install]
WantedBy=multi-user.target
